# provision/public/root/etc/profile.d/docker_ip.sh
docker_ip ()
{
  docker inspect -f "{{ .NetworkSettings.IPAddress }}" $1
}
