# provision/public/root/etc/profile.d/keychain.sh
eval $(keychain --eval --ignore-missing ~/.ssh/* 2> /dev/null)
