# provision/public/root/etc/profile.d/xterm-256color.sh
# BUGFIX: xterm color https://bbs.archlinux.org/viewtopic.php?id=175581
export TERM=xterm-256color
force_color_prompt=1
